# config/systemd/user/noctalia.service.d/override.conf
[Service]
Environment="QT_QPA_PLATFORMTHEME=gtk3"
